# regmik-erp-simple.service
[Unit]
Description=REGMIK ERP System
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/runner/workspace
ExecStart=/usr/bin/npx tsx server/production.ts
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=production
Environment=DATABASE_URL=postgresql://neondb_owner:npg_PQu4CAr9yIYq@ep-spring-flower-a552xsk9.us-east-2.aws.neon.tech/neondb?sslmode=require
Environment=SESSION_SECRET=regmik_erp_production_session_2025
Environment=PORT=3000

StandardOutput=journal
StandardError=journal
SyslogIdentifier=regmik-erp

[Install]
WantedBy=multi-user.target